# crontab
* * * * * /usr/local/bin/python /app/github/run.py >> /var/log/cron.log 2>&1
* * * * * /usr/local/bin/python /app/mongo/run.py >> /var/log/cron.log 2>&1
